# shieldbot-api.service
[Unit]
Description=ShieldAI Firewall API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=shieldbot
Group=shieldbot
WorkingDirectory=/opt/shieldbot
Environment="PATH=/opt/shieldbot/venv/bin"
EnvironmentFile=/opt/shieldbot/.env
ExecStart=/opt/shieldbot/venv/bin/uvicorn api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# Systemd hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/shieldbot
ProtectHome=true

[Install]
WantedBy=multi-user.target
