# shieldbot-api.service
[Unit]
Description=ShieldAI Firewall API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shieldbot
Environment="PATH=/opt/shieldbot/venv/bin"
EnvironmentFile=/opt/shieldbot/.env
ExecStart=/opt/shieldbot/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
